# daemon/session-nanny.service.in
[Unit]
Description=session-nanny

[Service]
BusName=net.sardemff7.SessionNanny
ExecStart=@libexecdir@/session-nanny
